# burp-client.timer
[Unit]
Description=Hourly runs burp timer

[Timer]
OnCalendar=hourly
RandomizedDelaySec=600
Persistent=true

[Install]
WantedBy=timers.target
